# root/etc/services.d/esphome/run
#!/usr/bin/with-contenv bash

echo "Starting ESPHome."
exec s6-setuidgid abc \
	esphome dashboard /config
