# root/etc/s6-overlay/s6-rc.d/init-config-esphome/run
#!/usr/bin/with-contenv bash

# make folderes
mkdir -p "${PIO_CACHE_BASE}"

# permissions
chown -R abc:abc \
    /config \
    /piolibs
